# modules/msm_backup_restore/files/cron.d/msm-backup-working-files-to-s3
0 * * * *   minecraft    /usr/local/venv/bin/python /usr/local/bin/msm-backup-working-files-to-s3.py
