# modules/msm_backup_restore/files/cron.d/msm-track-active-players
* * * * *   minecraft    /usr/local/venv/bin/python /usr/local/bin/msm-track-active-players.py
